# compose/aws/postgres-backup/crontab
# Daily pg_dumpall -> gzip -> S3. Runs as root (cron.d default) since this
# container has no other users configured.
0 3 * * * root /usr/local/bin/backup.sh >> /var/log/cron.log 2>&1
